# HungNV/Ghi_chep_python/Tools/tool-check-acepted_ssh/sshalert.service
[Unit]
Description=Tutorialinux sshalert service - alerts when someone logs in or sudos

[Install]
WantedBy=multi-user.target

[Service]
Type=simple
ExecStart=
Restart=always
